# systemd/pseudo-codex-dispatcher.service.d/20-browser-pool.conf
[Service]
Environment=PSEUDO_CODEX_MAX_WORKERS=2
